# backend/supervisord.conf
; =============================================================================
; supervisord.conf — TDS Sentinel API
; Process supervisor: auto-reinicia Flask si muere → elimina el 502
;
; Uso:
;   supervisord -c supervisord.conf          → arrancar (background)
;   supervisorctl -c supervisord.conf status → ver estado
;   supervisorctl -c supervisord.conf restart sentinel → reiniciar
;   supervisorctl -c supervisord.conf tail sentinel    → ver logs
; =============================================================================

[supervisord]
nodaemon          = false                       ; corre en background (daemon)
logfile           = /tmp/supervisord.log
logfile_maxbytes  = 5MB
logfile_backups   = 2
loglevel          = info
pidfile           = /tmp/supervisord.pid

; ── Control via socket (supervisorctl) ─────────────────────────────────────
[unix_http_server]
file = /tmp/supervisor.sock

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── Proceso principal: Flask HTTP ──────────────────────────────────────────
; HTTP puro — el proxy de Codespaces hace la terminación SSL externamente.
; No se necesitan certificados. Sin dependencia de openssl ni cryptography.
[program:sentinel]
command           = /workspaces/tds-sentinel-mvp-v3-2-0/backend/.venv/bin/python
                    /workspaces/tds-sentinel-mvp-v3-2-0/backend/server.py
directory         = /workspaces/tds-sentinel-mvp-v3-2-0/backend
environment       = PORT="5000"

; ── Restart automático ──────────────────────────────────────────────────────
autostart         = true    ; arranca al iniciar supervisord
autorestart       = true    ; reinicia si muere (sea cual sea el exit code)
startretries      = 10      ; intenta hasta 10 veces antes de marcar FATAL
startsecs         = 3       ; segundos que debe vivir para considerarse "up"
stopwaitsecs      = 10      ; tiempo para parar limpiamente

; ── Logs del proceso ────────────────────────────────────────────────────────
stdout_logfile          = /tmp/sentinel.log
stderr_logfile          = /tmp/sentinel.log
stdout_logfile_maxbytes = 5MB
stderr_logfile_maxbytes = 5MB
stdout_logfile_backups  = 2
stderr_logfile_backups  = 2
